# server/supervisord.conf
[program:node]
command = node cross-copy.js
directory = /home/dotcloud/current